# src-deb/1-second-stage/tmpfs-varlog-mysql/DEBIAN/postinst
#!/bin/sh

update-rc.d mysql-tmpfs defaults 90 10
