# src-deb/tmpfs-varlog-apache2/DEBIAN/postinst
#!/bin/sh

update-rc.d apache2-tmpfs defaults 90 10
